# scripts/systemd/avail-backup-verify.timer
[Unit]
Description=Weekly timer for AvailAI database backup verification (P1.5)
Documentation=https://github.com/TRIOSCS/Avail-AI-Test

[Timer]
# Sundays at 04:00 server time — off-peak, well clear of the 6-hourly
# db-backup cron (backup-cron.sh) so it always verifies a settled backup.
OnCalendar=Sun *-*-* 04:00:00
Persistent=true
Unit=avail-backup-verify.service

[Install]
WantedBy=timers.target
