# deploy/avail-xvfb.service
[Unit]
Description=Xvfb Virtual Display for AVAIL NC Worker
Documentation=https://github.com/TRIOSCS/Avail-AI-Test

[Service]
Type=simple
ExecStart=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -ac
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
